# docker/test/scripts/supervisord.conf
[supervisord]
nodaemon=true
logfile=%(ENV_LOG_PATH)s/supervisord-base.log
childlogdir=%(ENV_LOG_PATH)s

[program:init]
directory=/home/mauiusr
command=/usr/bin/bash /home/mauiusr/init.sh
autostart=true
startretries=0
stdout_logfile=%(ENV_LOG_PATH)s/init.stdout.log
redirect_stderr=true
priority=1
startsecs=0
autorestart=false
environment=
    INIT_PWSH_SCRIPT="%(ENV_INIT_PWSH_SCRIPT)s",
    INIT_BASH_SCRIPT="%(ENV_INIT_BASH_SCRIPT)s"

[program:emulator]
directory=/home/mauiusr
command=/usr/bin/bash /home/mauiusr/emulator.sh
autostart=true
autorestart=false
startretries=0
stdout_logfile=%(ENV_LOG_PATH)s/emulator.stdout.log
redirect_stderr=true
priority=2
stopsignal=TERM
stopwaitsecs=30

[program:appium]
directory=/home/mauiusr
command=/usr/bin/bash /home/mauiusr/appium.sh
autostart=true
autorestart=true
startretries=3
stdout_logfile=%(ENV_LOG_PATH)s/appium.stdout.log
redirect_stderr=true
priority=3

[program:androidportforward]
directory=/home/mauiusr
command=/usr/bin/bash /home/mauiusr/androidportforward.sh
autostart=true
autorestart=true
stdout_logfile=%(ENV_LOG_PATH)s/androidportforward.stdout.log
redirect_stderr=true
priority=4
